# supervisor.conf
[program:restapi]
command = /app/rest-api