# supervisor.conf
[program:restapi]
command = /app/rest-api
restart_when_binary_changed = true
[inet_http_server]
port = :9001
username=dalong
password=dalong