# supervisor.conf
[program:restapi]
command = /app/rest-api

[inet_http_server]
port = :9001